# etc/skel/.config/autostart/gemian-lock.desktop
[Desktop Entry]
Exec=xss-lock gemian-lock
Name=Screen Lock
Type=Application
Version=1.0
